# templates/profile/unison/server/unison.service.erb
[Unit]
Description=%i Unison
After=network.target

[Service]
Type=simple
RemainAfterExit=no
Restart=always
Environment=HOME=<%= @home %>

[Install]
WantedBy=multi-user.target
