# templates/profile/hathitrust/solr_lss/solr.service.erb
# Managed by puppet (nebula/profile/hathitrust/solr_lss/solr.service.erb)
[Unit]
Description=Solr for LSS
After=network.target htsolr-lss.mount

[Service]
User=solr
Type=forking
Environment="SOLR_INCLUDE=<%= @solr_home %>/solr.in.sh"
Environment="UMASK=002"
ExecStart=<%= @solr_bin %> start
ExecStop=<%= @solr_bin %> stop
Restart=always

[Install]
WantedBy=multi-user.target
