# templates/profile/fulcrum/fedora.service.erb
[Unit]
Description=Apache Tomcat running Fedora Repository

[Service]
User=fulcrum
Group=fulcrum
EnvironmentFile=-/etc/default/fedora
ExecStart=/usr/share/tomcat9/bin/catalina.sh run

[Install]
WantedBy=multi-user.target
