# templates/profile/fulcrum/fulcrum.target.erb
# Managed by puppet (nebula/profile/fulcrum/fulcrum.target.erb)
[Unit]
Description=Fulcrum
After=network-online.target

[Install]
WantedBy=multi-user.target

