# templates/profile/fulcrum/shibresponder.socket.erb
[Unit]
Description=Shibboleth Responder socket
After=network.target

[Socket]
ListenStream=127.0.0.1:9003
Accept=yes

[Install]
WantedBy=sockets.target
