# templates/profile/fulcrum/shibresponder.socket.erb
[Unit]
Description=Shibboleth Responder socket
After=network.target

[Socket]
ListenStream=/var/run/shibresponder.sock
SocketUser=_shibd
SocketGroup=_shibd
SocketMode=0666

[Install]
WantedBy=sockets.target
