# templates/named_instance/main.target.erb
# Managed by puppet (nebula/named_instance/main.target.erb)
[Unit]
Description=<%= @title -%> application
After=network.target
Requires=puma@<%= @title -%>.service
<% @subservices.each do |subservice| %>
Wants=<%= subservice -%>@<%= @title -%>.service
<% end %>

[Install]
WantedBy=multi-user.target
