# templates/named_instance/subservice.service.erb
# Managed by puppet (nebula/named_instance/subservice.service.erb)
[Unit]
Description=<%= @subservice -%> subservice for <%= @title %>
PartOf=<%= @title -%>.target
StartLimitIntervalSec=900
StartLimitBurst=5

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
UMask=0002
User=<%= @title %>
Group=<%= @title %>
Environment="RBENV_ROOT=<%= @rbenv_root -%>"
Environment="RAILS_ENV=production"
WorkingDirectory=<%= @path -%>/current
ExecStart=<%= @rbenv_root -%>/bin/rbenv exec bundle exec bin/<%= @subservice %>

TimeoutStartSec=60
Restart=on-failure
RestartSec=5

